# config/supervisor/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:apache]
command=/usr/sbin/apache2ctl -DFOREGROUND

[program:ssh]
command=/usr/sbin/sshd -D